# myapp.service

##generated by betterez automation

[Unit]
Description=simple node hello web server

[Service]
EnvironmentFile=/home/tal/programs/myapp.env
WorkingDirectory=/home/tal/programs/myapp
ExecStart=/usr/bin/npm start 
Restart=always

[Install]
WantedBy=multi-user.target

#systemctl daemon-reload
#systemctl start myapp.service
#sudo systemctl enable myapp.service
#build_number
